# playbooks/roles/mfl_api/templates/supervisor.conf
[program:gunicorn]
command=/opt/mfl_api_virtualenv/bin/gunicorn -w 3 config.wsgi:application
directory=/opt/mfl_api
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
