# playbooks/roles/api/templates/celery_worker.conf
;=================================
;celery worker supervisor config
; ==================================

; the name of your supervisord program
[program:mflcelery]

; Set full path to celery program if using virtualenv
command=/opt/mfl_api_virtualenv/bin/celery worker -A config --loglevel=INFO

; The directory to your Django project
directory=/opt/mfl_api

; If supervisord is run as the root user, switch users to this UNIX user account
; before doing any processing.
user={{mfl_celery_user}}

; Supervisor will start as many instances of this program as named by numprocs
numprocs=1

; Put process stdout output in this file
stdout_logfile=/var/log/celery/celery_worker.log

; Put process stderr output in this file
stderr_logfile=/var/log/celery/celery_worker.log

; If true, this program will start automatically when supervisord is started
autostart=true

; May be one of false, unexpected, or true. If false, the process will never
; be autorestarted. If unexpected, the process will be restart when the program
; exits with an exit code that is not one of the exit codes associated with this
; process’ configuration (see exitcodes). If true, the process will be
; unconditionally restarted when it exits, without regard to its exit code.
autorestart=true

; The total number of seconds which the program needs to stay running after
; a startup to consider the start successful.
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; if your broker is supervised, set its priority higher
; so it starts first
priority=998
environment=DEBUG="{{django_debug}}",EMAIL_HOST="{{email_host}}",EMAIL_HOST_USER="{{email_host_user}}",EMAIL_HOST_PASSWORD="{{email_host_password}}",REALTIME_INDEX="{{realtime_index}}",AWS_ACCESS_KEY_ID="{{storage_account}}",AWS_SECRET_ACCESS_KEY="{{storage_key}}",AWS_STORAGE_BUCKET_NAME="{{storage_bucket}}",AWS_DB_BACKUP_BUCKET="{{aws_db_backup_bucket}}"
